# scripts/systemd/catalog_server.service
[Unit]
Description=CCTools catalog server

[Service]
Type=forking
ExecStart=/opt/Yerba/scripts/catalog_server start 
Restart=on-abnormal
