# data/yerba
#!/bin/bash
### BEGIN INIT INFO
# Provides: yerba
# # Required-Start: $local_fs $syslog $remote_fs dbus
# # Required-Stop: $local_fs $syslog $remote_fs
# # Default-Start: 2 3 4 5
# # Default-Stop: 0 1 6
# # Short-Description: Start yerbad
# ### END INIT INFO

YERBA_ROOT=/opt/Yerba
CMD=$YERBA_ROOT/bin/yerbad
PROJECT_NAME=coge
CATALOG_HOST=localhost
CATALOG_PORT=1024
PYTHONPATH="/usr/local/lib/python2.7/site-packages:$YERBA_ROOT"

if [ "$(id -u)" != "0" ]; then
    echo "Script must be run as root." 1>&2
    exit 1
fi

start() {
    if [ "$(pgrep -f 'yerbad')" = "" ]; then
        sudo -u www-data env PROJECT_NAME=$PROJECT_NAME \
            CATALOG_HOST=$CATALOG_HOST CATALOG_PORT=$CATALOG_PORT \
            PYTHONPATH=$PYTHONPATH $CMD --queue-prefix $PROJECT_NAME \
            &> /dev/null &
        echo "Yerba started"
    else
        echo "Yerba already running"
    fi
}

stop() {
    echo "Yerba stopped"
    pkill -SIGKILL -f 'yerbad'
}

case "$1" in
    start)
        start
        exit 0
        ;;
    stop)
        stop
        exit 0
        ;;
    restart|reload|force-reload)
        stop
        start
        echo "Yerba restarted."
        exit 0
        ;;
    *)
        echo "Usage: $0 {start|stop|reload}" 1>&2
        exit 1
        ;;
esac
