# scripts/work_queue_pool.conf
# /etc/init/.conf

description "The cctools catalog server"
author "Evan Briones"

start on (local-filesystems and net-device-up IFACE=eth0)

stop on shutdown

respawn limit 30 60

pre-start script
    rm -f work_queue_pool.pid
end script

script
    CONFIG=/opt/Yerba/work_queue_pool/work_queue_pool.conf
    exec work_queue_pool -C localhost:1024 -A -c $CONFIG -d all
end script
